# ops/com.hawala.morning_brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Daily pre-trade signal brief → MACRO Telegram.
  Fires 09:20 IST Mon-Fri (after runner/VP-Trail signals would have fired
  at ~09:12-09:15, before the user is tempted to do anything).

  Install:
      cp ops/com.hawala.morning_brief.plist ~/Library/LaunchAgents/
      launchctl load ~/Library/LaunchAgents/com.hawala.morning_brief.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hawala.morning_brief</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/anaconda3/bin/python3</string>
        <string>-m</string>
        <string>ops.morning_brief</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/subhransubaboo/Claude Projects/Hawala v2/Hawala v2</string>

    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>20</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>20</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>20</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>20</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>20</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/subhransubaboo/Claude Projects/Hawala v2/Hawala v2/logs/macro_bot/morning_brief.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/subhransubaboo/Claude Projects/Hawala v2/Hawala v2/logs/macro_bot/morning_brief.err</string>
</dict>
</plist>
