# ops/com.hawala.autoheal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Self-healing routine. Fires 06:55 IST weekdays — BEFORE healthcheck (07:25)
  and the daily report (07:32). Runs ops/autoheal.py once: detects stale
  caches / broken fetchers, AUTO-FIXES what is safe (re-runs the fetcher),
  re-verifies, and escalates the rest to the MACRO Telegram channel.

  Runs once and exits (not a daemon) — so RunAtLoad=false, no KeepAlive.

  Install:
      cp ops/com.hawala.autoheal.plist ~/Library/LaunchAgents/
      launchctl load ~/Library/LaunchAgents/com.hawala.autoheal.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hawala.autoheal</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-i</string>
        <string>/opt/anaconda3/bin/python3</string>
        <string>ops/autoheal.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/subhransubaboo/Claude Projects/Hawala v2/Hawala v2</string>

    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>55</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>55</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>55</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>55</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>55</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/subhransubaboo/Claude Projects/Hawala v2/Hawala v2/logs/trade_bot/autoheal_launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/subhransubaboo/Claude Projects/Hawala v2/Hawala v2/logs/trade_bot/autoheal_launchd.err</string>
</dict>
</plist>
